# rootfs/etc/profile.d/locale.sh
#!/bin/sh
# Language settings
export LANG="C.UTF8"
export LC_ALL="C.UTF8"
